# config/aptrust_demo.conf
;; supervisord config file for APTrust processing apps
;; This file should be copied to /etc/supervisor/conf.d/aptrust.conf
;;
;; In addition to the programs listed below, there are two processes
;; that query the Processed Items list in Fluctus.
;;
;; bagman/apps/bucket_reader scans partner receiving buckets for bags
;; that have not yet been processed. It sends results into the
;; prepare_topic in NSQ.
;;
;; bagman/apps/request_reader checks Flucuts for requests to restore
;; intellectual objects and to delete generic files. It passes results
;; into the NSQ topics restore_topic and delete_topic.
;;
;; Both the bucket reader and the request reader run as cron jobs.
;; See the crontab entries for user ubuntu.
;;
;; -------------------------------------------------------------------------


;; Start the nsqd server first
[program:nsqd]
directory=/home/ubuntu/go/src/github.com/APTrust/bagman/bin
command=/home/ubuntu/go/src/github.com/APTrust/bagman/bin/apt_nsq_service -config /home/ubuntu/go/src/github.com/APTrust/bagman/bin/nsqd.demo.config
environment=GOPATH="/home/ubuntu/go",BAGMAN_HOME="/home/ubuntu/go/src/github.com/APTrust/bagman",PATH="$PATH:/home/ubuntu/go/bin:/usr/local/go/bin"
user=ubuntu
stopasgroup=true
stdout_logfile=/mnt/apt/logs/nsqd.stdout.log
stdout_logfile_maxbytes=1000000
stderr_logfile=/mnt/apt/logs/nsqd.stderr.log
stderr_logfile_maxbytes=1000000
stderr_logfile_backups=0
startsecs=10
stopwaitsecs=10
priority=100

;; apt_prepare reads from the NSQ prepare_channel, fetches new
;; bags from S3 receiving buckets, untars and validates them.
;; It sends info about valid bags to the NSQ store topic
[program:apt_prepare]
directory=/home/ubuntu/go/src/github.com/APTrust/bagman/bin
command=/home/ubuntu/go/src/github.com/APTrust/bagman/bin/apt_prepare -config=demo -env=/home/ubuntu/.bash_profile
environment=GOPATH="/home/ubuntu/go",BAGMAN_HOME="/home/ubuntu/go/src/github.com/APTrust/bagman"
user=ubuntu
stderr_logfile=/mnt/apt/logs/apt_prepare.stderr.log
startsecs=20
stopwaitsecs=30
priority=110

;; apt_store reads from the NSQ store_channel. It copies files that
;; were untarred by apt_prepare into the S3 preservation bucket,
;; then it sends results to the NSQ metadata_topic.
[program:apt_store]
directory=/home/ubuntu/go/src/github.com/APTrust/bagman/bin
command=/home/ubuntu/go/src/github.com/APTrust/bagman/bin/apt_store -config=demo -env=/home/ubuntu/.bash_profile
environment=GOPATH="/home/ubuntu/go",BAGMAN_HOME="/home/ubuntu/go/src/github.com/APTrust/bagman"
user=ubuntu
stderr_logfile=/mnt/apt/logs/apt_store.stderr.log
startsecs=20
stopwaitsecs=10
priority=120

;; apt_record reads from the metadata_channel in NSQ.
;; It records metadata about intellectual objects, generic files
;; and premis events in Fluctus.
[program:apt_record]
directory=/home/ubuntu/go/src/github.com/APTrust/bagman/bin
command=/home/ubuntu/go/src/github.com/APTrust/bagman/bin/apt_record -config=demo -env=/home/ubuntu/.bash_profile
environment=GOPATH="/home/ubuntu/go",BAGMAN_HOME="/home/ubuntu/go/src/github.com/APTrust/bagman"
user=ubuntu
stderr_logfile=/mnt/apt/logs/apt_record.stderr.log
startsecs=20
stopwaitsecs=10
priority=130


;; apt_trouble reads from the trouble_channel in NSQ. This topic/channel
;; contains information about bags whose ingest has failed for any reason
;; other than invalid file format and failed bag validation.
[program:apt_trouble]
directory=/home/ubuntu/go/src/github.com/APTrust/bagman/bin
command=/home/ubuntu/go/src/github.com/APTrust/bagman/bin/apt_trouble -config=demo -env=/home/ubuntu/.bash_profile
environment=GOPATH="/home/ubuntu/go",BAGMAN_HOME="/home/ubuntu/go/src/github.com/APTrust/bagman"
user=ubuntu
stderr_logfile=/mnt/apt/logs/apt_trouble.stderr.log
startsecs=10
stopwaitsecs=10
priority=140


;; apt_cleanup reads from the cleanup_channel in NSQ. It deletes the original
;; tar files from the partner receiving buckets after a bag has been successfully
;; ingested.
[program:apt_bag_delete]
directory=/home/ubuntu/go/src/github.com/APTrust/bagman/bin
command=/home/ubuntu/go/src/github.com/APTrust/bagman/bin/apt_bag_delete -config=demo -env=/home/ubuntu/.bash_profile
environment=GOPATH="/home/ubuntu/go",BAGMAN_HOME="/home/ubuntu/go/src/github.com/APTrust/bagman"
user=ubuntu
stderr_logfile=/mnt/apt/logs/apt_bag_delete.stderr.log
startsecs=10
stopwaitsecs=10
priority=150


;; apt_restore reads from the restore_topic in NSQ. It restores bags that
;; users have requested.
[program:apt_restore]
directory=/home/ubuntu/go/src/github.com/APTrust/bagman/bin
command=/home/ubuntu/go/src/github.com/APTrust/bagman/bin/apt_restore -config=demo -env=/home/ubuntu/.bash_profile
environment=GOPATH="/home/ubuntu/go",BAGMAN_HOME="/home/ubuntu/go/src/github.com/APTrust/bagman"
user=ubuntu
stderr_logfile=/mnt/apt/logs/apt_restore.stderr.log
startsecs=10
stopwaitsecs=10
priority=160


;; apt_delete reads from the delete_topic in NSQ. It deletes generic files
;; from the preservation bucket on user request.
[program:apt_file_delete]
directory=/home/ubuntu/go/src/github.com/APTrust/bagman/bin
command=/home/ubuntu/go/src/github.com/APTrust/bagman/bin/apt_file_delete -config=demo -env=/home/ubuntu/.bash_profile
environment=GOPATH="/home/ubuntu/go",BAGMAN_HOME="/home/ubuntu/go/src/github.com/APTrust/bagman"
user=ubuntu
stderr_logfile=/mnt/apt/logs/apt_file_delete.stderr.log
startsecs=10
stopwaitsecs=10
priority=170

;; apt_replicate reads from the replication_channel in NSQ.
;; It copies generic files from the preservation bucket
;; in Virginia to the replication bucket in Oregon.
;; *** THIS PROCESS SHOULD RUN ON THE RESTORE SERVER, IF POSSIBLE ***
[program:apt_replicate]
directory=/home/ubuntu/go/src/github.com/APTrust/bagman/bin
command=/home/ubuntu/go/src/github.com/APTrust/bagman/bin/apt_replicate -config=demo -env=/home/ubuntu/.bash_profile
environment=GOPATH="/home/ubuntu/go",BAGMAN_HOME="/home/ubuntu/go/src/github.com/APTrust/bagman"
user=ubuntu
stderr_logfile=/mnt/apt/logs/apt_replicate.stderr.log
startsecs=20
stopwaitsecs=10
priority=180

;; apt_failed_replication dumps data into a JSON file
;; describing failed attempts to copy a generic file to the
;; replication bucket in Oregon.
[program:apt_failed_replication]
directory=/home/ubuntu/go/src/github.com/APTrust/bagman/bin
command=/home/ubuntu/go/src/github.com/APTrust/bagman/bin/apt_failed_replication -config=production -env=/home/ubuntu/.bash_profile
environment=GOPATH="/home/ubuntu/go",BAGMAN_HOME="/home/ubuntu/go/src/github.com/APTrust/bagman"
user=ubuntu
stderr_logfile=/mnt/apt/logs/apt_failed_replication.stderr.log
startsecs=20
stopwaitsecs=10
priority=190

;; apt_fixity reads from the fixity_topic in NSQ. It runs fixity checks
;; on generic files in the S3 preservation bucket, comparing the SHA256
;; checksum in Fedora with the checksum it calculates from the actual
;; S3 file.
[program:apt_fixity]
directory=/home/ubuntu/go/src/github.com/APTrust/bagman/bin
command=/home/ubuntu/go/src/github.com/APTrust/bagman/bin/apt_fixity -config=production -env=/home/ubuntu/.bash_profile
environment=GOPATH="/home/ubuntu/go",BAGMAN_HOME="/home/ubuntu/go/src/github.com/APTrust/bagman"
user=ubuntu
stderr_logfile=/mnt/apt/logs/apt_fixity.stderr.log
startsecs=10
stopwaitsecs=10
priority=200


;; apt_failed_fixity reads from the failed_fixity_topic in NSQ.
;; It records information about fixity checks that could not be
;; completed, dumping JSON data to disk for human review.
[program:apt_failed_fixity]
directory=/home/ubuntu/go/src/github.com/APTrust/bagman/bin
command=/home/ubuntu/go/src/github.com/APTrust/bagman/bin/apt_failed_fixity -config=production -env=/home/ubuntu/.bash_profile
environment=GOPATH="/home/ubuntu/go",BAGMAN_HOME="/home/ubuntu/go/src/github.com/APTrust/bagman"
user=ubuntu
stderr_logfile=/mnt/apt/logs/apt_failed_fixity.stderr.log
startsecs=10
stopwaitsecs=10
priority=210

;; psdash monitoring listens on port 5000
[program:psdash]
command=/home/ubuntu/psdash/bin/psdash --agent
environment=PATH="/home/ubuntu/psdash/bin/:$PATH",PSDASH_CONFIG="/home/ubuntu/.psdash_config.py"
priority=220


;; Group all the worker programs, but allow nsqd to be started and stopped
;; separately.
[group:aptrust]
programs=apt_prepare,apt_store,apt_record,apt_trouble,apt_bag_delete,apt_restore,apt_file_delete,apt_replicate,apt_failed_replication,apt_fixity,apt_failed_fixity
