# config/dpn_services.conf
;; supervisord config file for DPN live ingest server
;; This file should be copied to /etc/supervisor/conf.d/dpn_services.conf
;;
;; In addition to the programs listed below, two other processes run as
;; cron jobs.
;;
;; 1) dpn_check_requests scans for replication requests and for APTrust
;;    bags that have been marked for ingest to DPN and sends results into
;;    the following topics in NSQ:
;;
;;       - APTrust bags bound for DPN go into dpn_package_topic
;;       - Bags to be replicated from other nodes go into dpn_copy_channel
;;
;; 2) dpn_sync connects to the REST services of other DPN nodes and copies
;;    DPN bag entries, replication requests and restore requests to our
;;    local DPN registry.
;;
;; See the crontab entries for user ubuntu.
;;
;; -------------------------------------------------------------------------


;; dpn_copy copies bags from remote nodes via rsync to our local
;; staging area. These are bags ingested elsewhere that APTrust
;; will be replicating
[program:dpn_copy]
directory=/home/ubuntu/go/src/github.com/APTrust/bagman/bin
command=/home/ubuntu/go/src/github.com/APTrust/bagman/bin/dpn_copy -config=production -env=/home/ubuntu/.bash_profile
environment=GOPATH="/home/ubuntu/go",BAGMAN_HOME="/home/ubuntu/go/src/github.com/APTrust/bagman"
user=ubuntu
stderr_logfile=/mnt/apt/logs/dpn_copy.stderr.log
startsecs=20
stopwaitsecs=30
priority=110

;; dpn_package packages locally ingested APTrust bags in DPN
;; format. Any APTrust bags marked for deposit to DPN have to
;; be packaged before being stored.
[program:dpn_package]
directory=/home/ubuntu/go/src/github.com/APTrust/bagman/bin
command=/home/ubuntu/go/src/github.com/APTrust/bagman/bin/dpn_package -config=production -env=/home/ubuntu/.bash_profile
environment=GOPATH="/home/ubuntu/go",BAGMAN_HOME="/home/ubuntu/go/src/github.com/APTrust/bagman"
user=ubuntu
stderr_logfile=/mnt/apt/logs/dpn_package.stderr.log
startsecs=20
stopwaitsecs=10
priority=120

;; dpn_validate validates bags copied from other nodes
;; before we send them off for storage.
[program:dpn_validate]
directory=/home/ubuntu/go/src/github.com/APTrust/bagman/bin
command=/home/ubuntu/go/src/github.com/APTrust/bagman/bin/dpn_validate -config=production -env=/home/ubuntu/.bash_profile
environment=GOPATH="/home/ubuntu/go",BAGMAN_HOME="/home/ubuntu/go/src/github.com/APTrust/bagman"
user=ubuntu
stderr_logfile=/mnt/apt/logs/dpn_validate.stderr.log
startsecs=20
stopwaitsecs=10
priority=130

;; dpn_store stores DPN bags in AWS S3 or Glacier.
[program:dpn_store]
directory=/home/ubuntu/go/src/github.com/APTrust/bagman/bin
command=/home/ubuntu/go/src/github.com/APTrust/bagman/bin/dpn_store -config=production -env=/home/ubuntu/.bash_profile
environment=GOPATH="/home/ubuntu/go",BAGMAN_HOME="/home/ubuntu/go/src/github.com/APTrust/bagman"
user=ubuntu
stderr_logfile=/mnt/apt/logs/dpn_store.stderr.log
startsecs=20
stopwaitsecs=10
priority=140

;; dpn_record records metadata about bags that we have
;; ingested or replicated. We record metadata about replicated
;; bags in the remote DPN (REST) registry that issued the
;; replication request. For bags going from APTrust to
;; DPN, we record data in both our own local DPN registry
;; and in Fluctus, so APTrust has a record of the item going
;; to DPN.
[program:dpn_record]
directory=/home/ubuntu/go/src/github.com/APTrust/bagman/bin
command=/home/ubuntu/go/src/github.com/APTrust/bagman/bin/dpn_record -config=production -env=/home/ubuntu/.bash_profile
environment=GOPATH="/home/ubuntu/go",BAGMAN_HOME="/home/ubuntu/go/src/github.com/APTrust/bagman"
user=ubuntu
stderr_logfile=/mnt/apt/logs/dpn_record.stderr.log
startsecs=20
stopwaitsecs=10
priority=150

;; dpn_trouble reads from the dpn_trouble_topic in NSQ. This topic/channel
;; contains information about DPN bags that could not be ingested.
[program:dpn_trouble]
directory=/home/ubuntu/go/src/github.com/APTrust/bagman/bin
command=/home/ubuntu/go/src/github.com/APTrust/bagman/bin/dpn_trouble -config=production -env=/home/ubuntu/.bash_profile
environment=GOPATH="/home/ubuntu/go",BAGMAN_HOME="/home/ubuntu/go/src/github.com/APTrust/bagman"
user=ubuntu
stderr_logfile=/mnt/apt/logs/dpn_trouble.stderr.log
startsecs=10
stopwaitsecs=10
priority=160

;; psdash monitoring listens on port 5000
[program:psdash]
command=/home/ubuntu/psdash/bin/psdash --agent
environment=PATH="/home/ubuntu/psdash/bin/:$PATH",PSDASH_CONFIG="/home/ubuntu/.psdash_config.py"
priority=200

;; Group all the DPN worker programs, so we can stop/start them
;; en mass.
[group:dpn]
programs=dpn_copy,dpn_package,dpn_validate,dpn_store,dpn_record,dpn_trouble
